# blue-value/cron.d/blue-value
## Run the Blue Value deployment schedule

* * * * * root (/algm/algm-lightsail-scripts/blue-value/deploy dev 2>&1 | logger)
