# blue-value/cron.d/blue-value
## Run the Blue Value deployment schedule

1 0 * * 0 root (/algm/algm-lightsail-scripts/blue-value/deploy dev 2>&1 | logger)
